# advance/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash /scripts/apache-start.sh
startretries=0

[program:redis]
command=/bin/bash /scripts/redis-start.sh
startretries=0
